# systemd/sports-live-poller.service
[Unit]
Description=Sports Data Live Score Poller
After=network.target

[Service]
Type=simple
User=sports
Group=sports
WorkingDirectory=/app
Environment=PYTHONPATH=/app:/app/dependencies/nba_api/src:/app/dependencies/mlb-statsapi
ExecStart=/usr/local/bin/python3 /app/scripts/poll_live_scores.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
