# sports-data-service/systemd/sports-live-poller.timer
[Unit]
Description=Sports Data Live Score Poller Timer
Requires=sports-live-poller.service

[Timer]
# Run every minute
# The service itself will decide whether to actually poll based on:
# 1. Polling hours (configured in settings)
# 2. Game states (poll if games in_progress or all upcoming)
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
