# ansible/templates/la-vuelta-refresh.service.j2
[Unit]
Description=La Vuelta Bundle Refresher ({{ deployment_env }})
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ansible
WorkingDirectory={{ deployment_dir }}/lavuelta-scraper
ExecStart={{ deployment_dir }}/lavuelta-scraper/.venv/bin/python {{ deployment_dir }}/lavuelta-scraper/refresh_if_due.py --outdir {{ deployment_dir }}/lavuelta-scraper
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
