# ansible/templates/sports-live-poller.service.j2
[Unit]
Description=Sports Data Live Score Poller ({{ deployment_env }})
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/docker-compose exec -T sports-service python scripts/poll_live_scores.py
WorkingDirectory={{ deployment_dir }}/sports-data-service
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

