# ansible/templates/sports-schedule-update.service.j2
[Unit]
Description=Sports Data Schedule Updater ({{ deployment_env }})
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/docker-compose exec -T sports-service python scripts/update_schedules.py
WorkingDirectory={{ deployment_dir }}/sports-data-service
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

